# packaging/debian/postinst
#!/usr/bin/env bash
set -euo pipefail

install -d -m 0750 -o root -g root /etc/ktc-mail /var/lib/ktc-mail /etc/ktc-mail/dkim
systemctl daemon-reload || true
systemctl enable ktc-mail-setup.service ktc-mail-firewall-monitor.timer ktc-mail-acme-renew.timer ktc-mail-backup.timer ktc-mail-exporter.timer ktc-mail-rate-limit.service || true
systemctl restart ktc-mail-setup.service || true

# Apply default SSH policy (key-only) if no custom config exists
if [[ ! -f /etc/ssh/sshd_config.d/ktc-mail.conf ]]; then
  /usr/lib/ktc-mail/ssh_policy.py apply 2>/dev/null || true
fi

# Load AppArmor profiles (conditional — AppArmor may not be present)
if command -v aa-enforce &>/dev/null; then
    for prof in /etc/apparmor.d/usr.lib.ktc-mail.*; do
        if [[ -f "$prof" ]]; then
            apparmor_parser -r "$prof" 2>/dev/null || true
            aa-enforce "$(basename "$prof")" 2>/dev/null || true
        fi
    done
elif command -v aa-complain &>/dev/null; then
    # If only complain mode tools are available, load in complain mode
    for prof in /etc/apparmor.d/usr.lib.ktc-mail.*; do
        if [[ -f "$prof" ]]; then
            apparmor_parser -r "$prof" 2>/dev/null || true
        fi
    done
fi

# Recommend unattended-upgrades package
if ! dpkg -l unattended-upgrades &>/dev/null 2>&1; then
    echo ""
    echo "  Recommendation: Install unattended-upgrades for automatic security updates."
    echo "    apt-get install unattended-upgrades"
    echo "  KTC Mail provides a config at /etc/apt/apt.conf.d/50ktc-mail-unattended-upgrades"
    echo ""
fi

cat <<'MSG'

╔═════════════════════════════════════════════════════════╗
║  KTC Mail setup GUI is available at                    ║
║                                                         ║
║     http://SERVER_IP:8080                                ║
║                                                         ║
║  Complete the DNS and TLS wizard before enabling         ║
║  production mail flow.                                   ║
║                                                         ║
║  Three things you need:                                  ║
║    1. Your domain (e.g. example.com)                     ║
║    2. Your DNS provider API key                          ║
║    3. Your admin email address                           ║
║                                                         ║
║  Everything else is auto-detected.                       ║
╚═════════════════════════════════════════════════════════╝
MSG
