# systemd/ktc-mail-exporter.service
[Unit]
Description=KTC Mail Prometheus metrics collector
Documentation=https://github.com/Dvalin21/KTC_MAIL
After=network-online.target

[Service]
Type=oneshot
User=ktc-mail
Group=ktc-mail
ExecStart=/usr/bin/ktc-mail metrics collect
IOSchedulingClass=idle
CPUSchedulingPolicy=idle
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/ktc-mail

[Install]
WantedBy=multi-user.target
