# systemd/ktc-mail-backup.service
[Unit]
Description=KTC Mail backup (restic)
Documentation=https://github.com/Dvalin21/KTC_MAIL
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/ktc-mail backup now
User=root
Group=root
TimeoutStartSec=7200
RemainAfterExit=yes
Restart=on-failure
RestartSec=300
StartLimitBurst=3
NoNewPrivileges=true
ProtectHome=false
ProtectSystem=full
PrivateTmp=true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
