# systemd/ktc-mail-firewall-monitor.service
[Unit]
Description=KTC Mail firewall policy monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/lib/ktc-mail/firewall_monitor.py

[Install]
WantedBy=multi-user.target
