# systemd/ktc-mail-firewall-monitor.service
[Unit]
Description=KTC Mail firewall policy monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/lib/ktc-mail/firewall_monitor.py
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/etc/ktc-mail /var/lib/ktc-mail /etc/nftables

[Install]
WantedBy=multi-user.target
