# systemd/ktc-mail-ops-check.timer
[Unit]
Description=Run KTC Mail ops checks hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
